# package/base-files/files/etc/init.d/example
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=99
start()
{
sleep 10
/usr/sbin/epoll_test&
}
stop()
{
killall epoll_test
}
